# tools/packages/systemd_services/oai-nr/nr-cuup.service
[Unit]
Description=nr-cuup

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/nr-cuup -O /etc/openairinterface/nr/gnb-cuup.sa.f1.conf
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
